# rpm/file-damage-analyzer.service
[Unit]
Description=File Damage Analyzer Service
After=network.target

[Service]
Type=simple
User=file-damage-analyzer
Group=file-damage-analyzer
ExecStart=/usr/bin/file-damage-analyzer --spring.config.location=file:/etc/file-damage-analyzer/application.properties
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target